# logging/scenarios/nginx_filebeat_es/supervisord.conf
# /etc/supervisor/conf.d/supervisord.conf 
[supervisord]
nodaemon=true

[program:filebeat]
command = service filebeat start 
[program:nginx]
command = service nginx start
